# btponto.crontab
# m h  dom mon dow   command
*/5 * * * * /usr/bin/btponto.py
